# data/fwupd-offline-update.service.in
[Unit]
Description=Updates device firmware whilst offline
Documentation=man:fwupdmgr
ConditionPathExists=@localstatedir@/lib/fwupd/offline-update

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=@libexecdir@/fwupd/fwupdoffline

[Install]
WantedBy=multi-user.target
